# conf.d/dosbox.conf
# this is the supervisor configuration, dosbox configuration
# is in file: dosbox-0.74-2.conf
[program:dosbox]
command=dosbox
autorestart=true
